# supervisor/event-consumer.conf
[program:eventconsumer]
command=/usr/bin/php "/opt/event-consumer/ec.php"
numprocs=1
process_name=event-consumer
autorestart=unexpected
exitcodes=0

[supervisord]
